# config/etc/cron.d/emailtodos
# at 5p PST email #todos
00 01 * * 2-6 django  cd ~/django_project && ./manage.py emailtodos 5
